# robocup-2022/field-computer-setup/ssl-game-controller.service
[Unit]
description=Launches SSL Game Controller Server

[Service]
Type=simple
ExecStart=%h/.local/bin/ssl-game-controller
Restart=always